# etc/supervisor/conf.d/zk_dockergen.conf
[program:zk_dockergen]
command=docker-gen -config /etc/dockergen/zookeeper.dockergen
directory=/etc/dockergen
autorestart=true
autostart=true
