# etc/supervisor/conf.d/zookeeper.conf
[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
user=zookeeper
autostart=false
autorestart=true
stopsignal=KILL
stopasgroup=true
